# docker/images/docker-spark/master/crontab
# The final argument number should match number of executor instances in the livy.conf file
0 20 * * * root /clueso/bin/table-compactor.sh /clueso/conf/application.conf local 5
